# contrib/systemd/fediwatcher.timer
[Unit]
Description=Fediwatcher timer (each 15 minutes)

[Timer]
OnBootSec=5min
OnCalendar=*:0/15

[Install]
WantedBy=basic.target
